# installer/config/systemd.service
[Unit]
Description=kyuuyomeisai Payroll Management System
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=kyuuyomeisai
Group=kyuuyomeisai
WorkingDirectory=/opt/kyuuyomeisai/backend
Environment=PATH=/opt/kyuuyomeisai/venv/bin
EnvironmentFile=/opt/kyuuyomeisai/.env
ExecStart=/opt/kyuuyomeisai/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
